# elixir-phoenix/overlay/usr/lib/systemd/system/cog.service
[Unit]
Description=Cog Webkit Launcher
After=network-online.target ref-elixir.service
Wants=network-online.target ref-elixir.service

[Service]
# Ensure environment substitution works
EnvironmentFile=-/etc/default/cog-avocado

# Default URL (can be overridden in EnvironmentFile)
Environment=COG_URL=https://avocadolinux.org

# Wayland/DRM runtime directory
Environment=XDG_RUNTIME_DIR=/run/user/0

# Writable cache for GPU shader cache
Environment=HOME=/tmp

# Build the CA bundle before starting cog
ExecStartPre=/bin/sh -c "mkdir -p /run/user/0"
ExecStartPre=/bin/sh -c "cat /usr/share/ca-certificates/mozilla/*.crt > /etc/ssl/certs/ca-certificates.crt"

# Launch cog
ExecStart=/usr/bin/cog --platform-params='renderer=gles' ${COG_URL}

Restart=always

[Install]
WantedBy=multi-user.target
