# iphone-travel-router/app/overlay/usr/lib/systemd/system/iphone-pair.service
[Unit]
Description=Pair attached iPhone via libimobiledevice
After=usbmuxd.service
Wants=usbmuxd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/iphone-pair.sh
StandardOutput=journal
StandardError=journal
