# rzv2n-drpai-yolo/app/overlay/usr/lib/systemd/system/weston.service.d/avocado-rzv2n-home.conf
# Avocado's rootfs is read-only erofs and ships an empty /home, so the
# upstream weston.service's WorkingDirectory=/home/weston fails at chdir
# (status=200/CHDIR) and the compositor never starts.
#
# Move weston's working directory to /var/lib/weston (writable btrfs).
# StateDirectory= makes systemd create the dir on each start with
# 0700 weston:weston ownership, which is what weston-autologin / pam
# expect for HOME-style behavior.

[Service]
StateDirectory=weston
WorkingDirectory=/var/lib/weston
Environment=HOME=/var/lib/weston
Environment=XDG_RUNTIME_DIR=/run/user/970
